# launchd_plist/com.myaiemployee.news-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AI 每日情报 — 独立 one-shot LaunchAgent

  - RunAtLoad + StartInterval=3600：登录后刷新一次，此后约每小时刷新公开 Feed
  - KeepAlive=false：来源失败不触发 launchd 热循环
  - ProcessType=Standard：禁止 Background（节流会把 1h 拖到 2–3h，撞 P3 news_run_gap）
  - 不替代 Dashboard；不执行业务服务启停、不访问 IMAP/Keychain、不发送 SMTP
  - 系统睡眠仍会冻结 StartInterval；P3 长稳需机器保持醒着或另行 caffeinate
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myaiemployee.news-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/$USER/bin/my-ai-employee-news-refresh-runner</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/$USER/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/$USER</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>KeepAlive</key>
    <false/>

    <key>ProcessType</key>
    <string>Standard</string>

    <key>StandardOutPath</key>
    <string>/Users/$USER/Library/Logs/MyAIEmployee/news-refresh.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/$USER/Library/Logs/MyAIEmployee/news-refresh.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/$USER</string>
</dict>
</plist>
