# launchd_plist/com.myaiemployee.digital-employee.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Day 2 launchd - 菜单栏 + Dashboard 开机自启

  wrapper ~/bin/my-ai-employee-start 调 ops/start-digital-employee.sh start。
  RunAtLoad=true:登录后自动拉起(撞坑 #81:TCC 须先授权 Python.framework 3.12)。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myaiemployee.digital-employee</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/$USER/bin/my-ai-employee-start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/$USER/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/$USER</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/$USER/Library/Logs/MyAIEmployee/digital-employee.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/$USER/Library/Logs/MyAIEmployee/digital-employee.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/$USER</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
