# launchd_plist/com.myaiemployee.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  D10.3 launchd plist for v0.1 数字生活月报 + launchd 保活

  触发时机:每月 1 号 09:00 触发 monthly_report generate(沿 week2-mvp.md L224 + 审计员触发)
  部署位置:~/Library/LaunchAgents/com.myaiemployee.agent.plist(用户级)
  装载方式:launchctl load -w plist(沿 launchd_install.sh 范本)

  设计决策(2026-06-15 锁定):
    RunAtLoad=false(避免开机就触发;用户期望是每月 1 号 09:00)
    StandardOutPath/StandardErrorPath 走 ~/Library/Logs/MyAIEmployee/(沿 macOS 标准日志位置)
    StartCalendarInterval:Month=1, Day=1, Hour=9, Minute=0(每月 1 号 09:00 准点)
    ProgramArguments 沿 ~/bin/ 部署(launchd 在 macOS Sequoia TCC 不允许 Documents/ 路径)

  验收:
    launchctl list | grep myaiemployee 必见 com.myaiemployee.agent
    (2026-06-15 D10.5.3 修正: 早期 docstring 写"launchctl print 必能读",
     但 macOS Sequoia/Tahoe launchctl print user/<uid>/<label> 域命令
     不稳定,会报 "Could not find service" / "Bad request",沿 Spike A §3.3
     教训改以"launchctl list 必见 label"为唯一真理源)
    日志文件可读:log show with predicate process = my-ai-employee-monthly-report

  ⚠️ 关于 ProgramArguments 里的 $(/bin/date '+%Y-%m')(2026-06-15 D10.5.3 P3-1):
    launchd 不会做 shell 展开,L31 的 $(/bin/date '+%Y-%m') 是字面字符串,
    实际 monthly_report.py 收到的 month 参数就是字面 "2026-06" 之类。
    v0.2 wrapper 脚本决策(把 monthly_report.py 包到一个 shell 脚本里,
    由 wrapper 自己 date 取月份再调 monthly_report.py generate)留待 v0.2 启动评估。
    当前 v0.1.0 范围:每月 1 号 09:00 触发时,需要手动同步更新 plist 的 month。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.myaiemployee.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/$USER/bin/my-ai-employee-monthly-report</string>
        <string>generate</string>
        <string>--month</string>
        <string>$(/bin/date '+%Y-%m')</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/$USER/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/$USER</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Month</key>
        <integer>1</integer>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/$USER/Library/Logs/MyAIEmployee/agent.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/$USER/Library/Logs/MyAIEmployee/agent.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/$USER</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
